# systemd/speedtest-monitor.timer
[Unit]
Description=Speedtest Monitor Timer
Documentation=https://github.com/SokolovMO/speedtest_monitor

[Timer]
# ============================================================================
# Timer Configuration Examples
# ============================================================================
# Uncomment ONE of the timer configurations below:

# --- Default: Every hour ---
OnCalendar=hourly

# --- Alternative Options (uncomment ONE to use instead) ---

# Interval-based (relative to last completion)
#OnBootSec=1min
#OnUnitActiveSec=1h
#AccuracySec=1min

# Every 30 minutes (:00 and :30)
#OnCalendar=*:0/30

# Every 10 minutes (:00, :10, :20, :30, :40, :50)
#OnCalendar=*:0/10

# Every 5 minutes
#OnCalendar=*:0/5

# Every minute
#OnCalendar=*:*

# Daily at midnight
#OnCalendar=*-*-* 00:00:00

# Every Monday at 9:00 AM
#OnCalendar=Mon *-*-* 09:00:00

# ============================================================================
# Documentation:
# - man systemd.timer
# - https://www.freedesktop.org/software/systemd/man/systemd.timer.html
# - https://www.freedesktop.org/software/systemd/man/systemd.time.html
#
# Test calendar syntax:
#   systemd-analyze calendar "hourly"
#   systemd-analyze calendar "*:0/10"
#
# After changes, run:
#   sudo systemctl daemon-reload
#   sudo systemctl restart speedtest-monitor.timer
#   systemctl list-timers speedtest-monitor.timer
# ============================================================================

[Install]
WantedBy=timers.target
